# .config/systemd/user/cabal-update.service
[Unit]
Description=Update cabal's package list

[Service]
Type=oneshot
ExecStart=cabal update
