# .config/sh/profile.d/cabal.sh
# Add binaries built with cabal to the path.
prepend_path "$HOME/.cabal/bin"
