# .config/sh/profile.d/firefox.sh
# shellcheck shell=bash

# Enable Wayland support
export MOZ_ENABLE_WAYLAND=1
