# .config/sh/profile.d/headphones_hissing.sh
# headphones_hissing.sh

# Sets the headphone microphone boost volume to 10dB to disable the hissing
# sound when using headphones.

#amixer -c PCH cset 'name=Headphone Mic Boost Volume' 1 >/dev/null
