# .config/sh/profile.d/sshaskpass.sh
# Let ssh-add obtain passwords from the user, even if not connected to a
# terminal, through ksshaskpass.
export SSH_ASKPASS='/usr/bin/ksshaskpass'
