# .config/sh/profile.d/tex.sh
# The editor that TeX starts to correct errors
export TEXEDIT="$EDITOR +%d %s"
